# services-user/miemo-record.service
[Unit]
Description=Miemo record service

[Service]
Type=simple
ExecStart=runMiemoRecord
Restart=always

[Install]
WantedBy=default.target
