# services-user/miemo-webapp.service
[Unit]
Description=Miemo webapp server

[Service]
Type=simple
ExecStart=runMiemoWebapp
Restart=always

[Install]
WantedBy=default.target
