# services-user/miemo-game.service
[Unit]
Description=Miemo game client

[Service]
Type=simple
ExecStart=runMiemoGame
Restart=always

[Install]
WantedBy=default.target
