# services-user/miemo-holo.service
[Unit]
Description=Miemo holo server

[Service]
Type=simple
ExecStart=runMiemoHolo
Restart=always

[Install]
WantedBy=default.target
