# services-user/miemo-audio.service
[Unit]
Description=Miemo Audio Manager server

[Service]
Type=simple
ExecStart=miemo-audio
Restart=always

[Install]
WantedBy=default.target
