# services-user/miemo-static.service
[Unit]
Description=Miemo static client

[Service]
Type=simple
ExecStart=runMiemoStatic
Restart=always

[Install]
WantedBy=default.target
