# services-user/miemo-webclient.service
[Unit]
Description=Miemo webapp client

[Service]
Type=idle
ExecStart=runMiemoWebClient
Restart=always

[Install]
WantedBy=default.target
